# scripts/server/desktop/per-user/wrappedhl.service
[Unit]
Description=Hyprland Wayland Compositor

[Service]
ExecStart=%h/.local/bin/wrappedhl
Restart=on-failure

[Install]
WantedBy=default.target

