# programs/_a-docker-starter/my-docker-flag.service
[Unit]
Description=Clear marker file in home directory on user startup after boot

[Service]
Type=oneshot
ExecStart=/bin/rm -f ~/.myDockerFlag

[Install]
WantedBy=default.target